# finAPI/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND

[program:uvicorn]
directory=/app
command=uvicorn stockapi.main:app
